# docker/services/sdrplay/run
#!/bin/bash

if [ "${radio,,}" != "sdrplay" ]; then
    s6-svc -O .
    exit 0
fi

exec /opt/sdrplay_api/sdrplay_apiService

