# docker/services/sdrplay/run
#!/command/execlineb -P
/opt/sdrplay_api/sdrplay_apiService
