# docker/services/sdrplay/run
#!/bin/bash

if [ "${RADIO,,}" != "sdrplay" ]; then
    echo "SDRplay radio not selected - not starting API service"
    s6-svc -O .
    exit 0
fi

exec /opt/sdrplay_api/sdrplay_apiService

